# deploy/daily-dryrun-review.service
[Unit]
Description=Trading Bot - Daily Dry-Run Review
After=network-online.target

[Service]
Type=oneshot
User=trading
Group=trading
WorkingDirectory=/opt/trading_2
EnvironmentFile=/opt/trading_2/.env
ExecStart=/bin/bash /opt/trading_2/scripts/run_daily_dryrun_review.sh
StandardOutput=append:/opt/trading_2/logs/daily_review.log
StandardError=append:/opt/trading_2/logs/daily_review.log
