# deploy/node-exporter-textfile.service
[Unit]
Description=Trading Bot - Link textfile metrics into Node Exporter collector directory
After=prometheus-node-exporter.service
Requires=scorecard-status.service

[Service]
Type=oneshot
User=root
# Symlink our prom file into the system Node Exporter textfile directory
ExecStart=/bin/ln -sfn /opt/trading_2/results/scorecards/textfile/trading_scorecard.prom /var/lib/prometheus/node-exporter/trading_scorecard.prom
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
