# deploy/scorecard.timer
[Unit]
Description=Trading Bot – Weekly Scorecard Timer
Requires=scorecard.service

[Timer]
# Run every Sunday at 09:00 local time
OnCalendar=Sun *-*-* 09:00:00
Persistent=true

[Install]
WantedBy=timers.target
