# deploy/tuning-cycle.timer
[Unit]
Description=Trading Bot - Weekly Tuning Cycle Timer
Requires=tuning-cycle.service

[Timer]
# Run every Sunday at 09:20 local time (after weekly scorecard at 09:00)
OnCalendar=Sun *-*-* 09:20:00
Persistent=true

[Install]
WantedBy=timers.target
